# service_files/openrc/project
#!/sbin/openrc-run
command="/usr/local/bin/block_daemon.sh"
description="Domain blocker tool"
